# systemd/pi-avatar-renderer.service
[Unit]
Description=Pi Avatar renderer
After=local-fs.target

[Service]
Type=simple
EnvironmentFile=-/etc/pi-avatar/avatar.env
Environment=XDG_RUNTIME_DIR=/run/pi-avatar
Environment=SDL_VIDEODRIVER=kmsdrm
WorkingDirectory=/opt/pi-avatar
ExecStartPre=/bin/sh -c '/usr/bin/setterm --cursor off --blank 0 --powerdown 0 > /dev/tty1'
ExecStart=/opt/pi-avatar/.venv/bin/python /opt/pi-avatar/renderer.py --config /etc/pi-avatar/avatar.yaml
ExecStopPost=/bin/sh -c '/usr/bin/setterm --cursor on > /dev/tty1'
Restart=on-failure
RestartSec=2
RuntimeDirectory=pi-avatar
RuntimeDirectoryMode=0700
StandardInput=tty
TTYPath=/dev/tty1

[Install]
WantedBy=multi-user.target
